# plists/com.user.battery-schedule.day.0900.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.battery-schedule.day.0900</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/davidchu/.battery-schedule/bin/apply.sh</string>
        <string>day</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Minute</key>
            <integer>0</integer>
            <key>Hour</key>
            <integer>9</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/davidchu/.battery-schedule/logs/schedule.day.0900.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/davidchu/.battery-schedule/logs/schedule.day.0900.log</string>
</dict>
</plist>